# deploy/systemd/stavki-v5.service
[Unit]
Description=Stavki Value Finder V5 (Pipeline Run)
After=network.target

[Service]
Type=oneshot
User=stavki
WorkingDirectory=/opt/stavki_value_system
EnvironmentFile=/opt/stavki_value_system/.env
ExecStart=/usr/bin/python3 scripts/run_scheduler.py --now --telegram

[Install]
WantedBy=multi-user.target
