# deploy/stavki-bot.service
[Unit]
Description=STAVKI Telegram Bot Interface
After=network.target

[Service]
# User=ubuntu  <-- UNCOMMENT AND SET CORRECT USER
# Group=ubuntu
WorkingDirectory=/home/ubuntu/stavki_value_system
EnvironmentFile=/home/ubuntu/stavki_value_system/.env
ExecStart=/usr/bin/python3 scripts/run_telegram_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
