# deploy/stavki-beat.service
[Unit]
Description=STAVKI Celery Beat (Scheduler)
After=network.target redis-server.service

[Service]
User=serni13678
Group=serni13678
WorkingDirectory=/home/serni13678/stavki_value_system
EnvironmentFile=/home/serni13678/stavki_value_system/.env
ExecStart=/home/serni13678/stavki_value_system/venv/bin/celery -A src.celery_app beat --loglevel=info
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
