# deploy/stavki-bot.service
[Unit]
Description=STAVKI Telegram Bot Interface
After=network.target

[Service]
User=serni13678
Group=serni13678
WorkingDirectory=/home/serni13678/stavki_value_system
EnvironmentFile=/home/serni13678/stavki_value_system/.env
ExecStart=/home/serni13678/stavki_value_system/venv/bin/python3 scripts/run_telegram_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
